# todoapp.service
[Unit]
Description=My To-Do App
[Service]
WorkingDirectory=/opt/todoapp
ExecStart="poetry run flask run -h 0.0.0.0 -p 80"